# config/supervisor/jake-bench.conf
[hover-jet]
file=/tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisord/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=error
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=root
childlogdir=/var/log/supervisord/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:logger]
environment=LOGGED_CHANNELS="camera_image_channel,fiducial_detection_channel,imu",
            LOG_BASE_PATH=/logs/
command=/jet/bin/run/message_logger_bq_main
stderr_logfile = /var/log/supervisord/logger-stderr.log
stdout_logfile = /var/log/supervisord/logger-stdout.log
stopsignal=INT

[program:imu]
command=/jet/bin/run/imu_balsaq_main
stderr_logfile = /var/log/supervisord/imu-stderr.log
stdout_logfile = /var/log/supervisord/imu-stdout.log
stopsignal=INT

[program:camera]
command=/jet/bin/run/camera_balsaq_main
stderr_logfile = /var/log/supervisord/camera_balsaq_main-stderr.log
stdout_logfile = /var/log/supervisord/camera_balsaq_main-stdout.log
stopsignal=INT

[program:fiducial_detection]
command=/jet/bin/run/fiducial_detection_balsaq_main
stderr_logfile = /var/log/supervisord/fiducial_detection_balsaq_main-stderr.log
stdout_logfile = /var/log/supervisord/fiducial_detection_balsaq_main-stdout.log
stopsignal=INT

[inet_http_server]
port=*:80
; username=
; password=
