# config/supervisor/joystick_host.conf
[hover-jet]
file=/tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisord/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=error
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=root
childlogdir=/var/log/supervisord/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:timesync_client]
command=/jet/bin/run/timesync_client_bq_main
stderr_logfile = /var/log/supervisord/timesync-client-stderr.log
stdout_logfile = /var/log/supervisord/timesync-client-stdout.log
stopsignal=INT

[program:turbine_control]
command=/jet/bin/run/turbine_control
stderr_logfile = /var/log/supervisord/turbine-control-stderr.log
stdout_logfile = /var/log/supervisord/turbine-control-stdout.log
stopsignal=INT

[program:pitch_roll_control]
command=/jet/bin/run/pitch_roll_control
stderr_logfile = /var/log/supervisord/pitch_roll-control-stderr.log
stdout_logfile = /var/log/supervisord/pitch_roll-control-stdout.log
stopsignal=INT

[program:yaw_control]
command=/jet/bin/run/yaw_control
stderr_logfile = /var/log/supervisord/yaw-control-stderr.log
stdout_logfile = /var/log/supervisord/yaw-control-stdout.log
stopsignal=INT

[inet_http_server]
port=*:80
; username=
; password=
