# config/supervisor/supervisor.conf
[hover-jet]
file=/tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisord/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=error
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=root
childlogdir=/var/log/supervisord/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:turbine_interface]
command=/jet/run/turbine_interface_bq_main
stderr_logfile = /var/log/supervisord/turbine_interface-stderr.log
stdout_logfile = /var/log/supervisord/turbine_interface-stdout.log

[program:logger]
command=/jet/run/message_logger_bq_main
stderr_logfile = /var/log/supervisord/logger-stderr.log
stdout_logfile = /var/log/supervisord/logger-stdout.log

[program:imu]
command=/jet/run/imu_balsaq_main
stderr_logfile = /var/log/supervisord/imu-stderr.log
stdout_logfile = /var/log/supervisord/imu-stdout.log

[program:servo_interface]
command=/jet/run/servo_balsaq_main
stderr_logfile = /var/log/supervisord/servo_interface-stderr.log
stdout_logfile = /var/log/supervisord/servo_interface-stdout.log

[inet_http_server]
port=*:80
; username=
; password=
